# deploy/seer-crawler-arbitrum-sepolia.service
[Unit]
Description=Seer crawler service for Arbitrum Sepolia blockchain
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
WorkingDirectory=/home/ubuntu/seer
EnvironmentFile=/home/ubuntu/seer-secrets/app.env
Restart=on-failure
RestartSec=15s
ExecStart=/home/ubuntu/seer/seer crawler --chain arbitrum_sepolia --confirmations 10 --threads 2
SyslogIdentifier=seer-crawler-arbitrum-sepolia

[Install]
WantedBy=multi-user.target
