# deploy/seer-synchronizer-ethereum.service
[Unit]
Description=Seer synchronizer service for Ethereum blockchain
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
WorkingDirectory=/home/ubuntu/seer
EnvironmentFile=/home/ubuntu/seer-secrets/app.env
Restart=on-failure
RestartSec=15s
ExecStart=/home/ubuntu/seer/seer synchronizer --chain ethereum
SyslogIdentifier=seer-synchronizer-ethereum

[Install]
WantedBy=multi-user.target
