# deploy/seer-synchronizer-sepolia.service
[Unit]
Description=Seer synchronizer service for Sepolia blockchain
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
WorkingDirectory=/home/ubuntu/seer
EnvironmentFile=/home/ubuntu/seer-secrets/app.env
Restart=on-failure
RestartSec=15s
ExecStart=/home/ubuntu/seer/seer synchronizer --chain sepolia
SyslogIdentifier=seer-synchronizer-sepolia

[Install]
WantedBy=multi-user.target
